# deploy/preview/ovh/setup/launcher.service
[Unit]
Description=Atrium OVH preview launcher
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=__ATRIUM_SERVICE_USER__
Group=__ATRIUM_SERVICE_USER__
SupplementaryGroups=docker
WorkingDirectory=__ATRIUM_REPO_ROOT__
EnvironmentFile=__ATRIUM_LAUNCHER_ENV__
Environment=PYTHONUNBUFFERED=1
# The docker CLI writes buildx state under $DOCKER_CONFIG (default ~/.docker).
# Keep it on a persistent, service-owned path so the warm buildx builder and its
# cache survive restarts (and so it never depends on $HOME being writable).
Environment=DOCKER_CONFIG=/var/lib/atrium-preview/.docker
# Same for helm's repo cache and kubectl's discovery cache: persist them rather
# than losing them with the tmpfs home. Without this helm fails with
# "no cached repository ... (try 'helm repo update')".
Environment=HELM_CACHE_HOME=/var/lib/atrium-preview/.helm/cache
Environment=HELM_CONFIG_HOME=/var/lib/atrium-preview/.helm/config
Environment=HELM_DATA_HOME=/var/lib/atrium-preview/.helm/data
Environment=KUBECACHEDIR=/var/lib/atrium-preview/.kube/cache
# Where previewctl drops per-preview vhost fragments. Must match the dir the
# shared Caddy imports; previewctl's built-in default (/etc/caddy/conf.d) is not
# this box's layout.
Environment=ATRIUM_PREVIEW_CADDY_CONF_DIR=__ATRIUM_CADDY_CONF_DIR__
# Must match the wildcard cert host in the shared Caddyfile.
Environment=ATRIUM_PREVIEW_DOMAIN=__ATRIUM_PREVIEW_DOMAIN__
ExecStart=/usr/bin/python3 __ATRIUM_REPO_ROOT__/deploy/preview/ovh/launcher.py
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# tmpfs, not read-only: this service shells out to docker, helm, kubectl, git and
# just, all of which write caches under $HOME. A read-only home makes each of them
# fail at its first write. A private tmpfs home keeps the sandbox tight while
# letting any tool scribble; the caches worth persisting are redirected above.
ProtectHome=tmpfs
# __ATRIUM_CADDY_CONF_DIR__ lives under /etc, which ProtectSystem=full mounts
# read-only — without this exception the route phase dies on "[Errno 30]
# Read-only file system".
ReadWritePaths=/var/lib/atrium-preview /var/cache/atrium-preview __ATRIUM_REPO_ROOT__/deploy/preview/ovh __ATRIUM_CADDY_CONF_DIR__

[Install]
WantedBy=multi-user.target
