# deploy/preview/ovh/setup/launcher.service
[Unit]
Description=Atrium OVH preview launcher
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=__ATRIUM_SERVICE_USER__
Group=__ATRIUM_SERVICE_USER__
SupplementaryGroups=docker
WorkingDirectory=__ATRIUM_REPO_ROOT__
EnvironmentFile=__ATRIUM_LAUNCHER_ENV__
Environment=PYTHONUNBUFFERED=1
# ProtectHome=read-only makes the service user's home unwritable, and the docker
# CLI writes buildx state under $DOCKER_CONFIG (default ~/.docker). Point it at a
# writable, service-owned dir or every `docker build` fails instantly.
Environment=DOCKER_CONFIG=/var/lib/atrium-preview/.docker
ExecStart=/usr/bin/python3 __ATRIUM_REPO_ROOT__/deploy/preview/ovh/launcher.py
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/atrium-preview /var/cache/atrium-preview __ATRIUM_REPO_ROOT__/deploy/preview/ovh

[Install]
WantedBy=multi-user.target
